# templates/com.niwa.executor.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key><string>com.niwa.executor</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{VENV_PYTHON}}</string>
        <string>-m</string>
        <string>app.executor</string>
    </array>
    <key>WorkingDirectory</key><string>{{REPO_DIR}}/backend</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>NIWA_CLAUDE_CLI</key><string>{{CLAUDE_CLI_PATH}}</string>
        <key>NIWA_CONFIG_PATH</key><string>{{HOME}}/.niwa/config.toml</string>
    </dict>
    <key>StandardOutPath</key><string>{{HOME}}/.niwa/logs/executor.log</string>
    <key>StandardErrorPath</key><string>{{HOME}}/.niwa/logs/executor.log</string>
    <key>KeepAlive</key><true/>
    <key>RunAtLoad</key><true/>
</dict>
</plist>
